# system/etc/sv/crond/run
#!/bin/sh
. /etc/runit/functions
[ -r conf ] && . ./conf
[ -d "$CRONDIR" ] || mkdir -p "$CRONDIR"
busybox crond -c "$CRONDIR" -f 2>&1
