# system/etc/sv/crond/run
#!/bin/sh
. /etc/runit/functions
[ -r conf ] && . ./conf
[ -d "$CRONDIR" ] || mkdir -p "$CRONDIR"
exec_cmd crond -c "$CRONDIR" -f
