# system/etc/sv/sshd/run
#!/data/data/com.termux/files/usr/bin/sh
# Termux Services - Rootful Sshd Service
# Prepare environment
. /etc/runit/functions
TERMUX_HOME=/data/data/com.termux/files/home
[ -r conf ] && . ./conf
TERMUX_EMPTY=$PREFIX/var/empty
# Built-in functions
generate_ssh_config() {
    cat > "$SSHD_CONFIG_FILE" <<- EOF
# Default sshd_config file generated by Termux Services
# Modify $SSHD_CONFIG_FILE or add files in $SSHD_CONFIG_DIR/ to customize
Port 8822
PrintMotd yes
PermitRootLogin yes
PasswordAuthentication no
AuthorizedKeysFile $AUTHORIZED_KEYS_FILE
Subsystem sftp $PREFIX/libexec/sftp-server
Include $SSHD_CONFIG_DIR/*.conf
EOF
}

# Generate ssh necessary files and directories
[ ! -d "$RUNDIR_EMPTY" ] && mkdir -p "$RUNDIR_EMPTY" && chmod 755 "$RUNDIR_EMPTY"

if [ ! -e "$SSH_CONFIG_DIR" ]; then
    mkdir -p "$SSHD_CONFIG_DIR"
    generate_ssh_config > $SSHD_CONFIG_FILE
fi

if [ -f "$TERMUX_HOME/$AUTHORIZED_KEYS_FILE" ] && [ ! -f "$ROOT_HOME/$AUTHORIZED_KEYS_FILE" ]; then
    install -Dvm 600 "$TERMUX_HOME/$AUTHORIZED_KEYS_FILE" "$ROOT_HOME/$AUTHORIZED_KEYS_FILE"
    chown root:root "$ROOT_HOME/$AUTHORIZED_KEYS_FILE"
fi

# Run sshd
if command -v sshd >/dev/null 2>&1; then
    exec unshare -m sh -c "
        mount --bind $RUNDIR_EMPTY $TERMUX_EMPTY
        mount --bind $ROOT_HOME $TERMUX_HOME
        exec sshd -D -f $SSHD_CONFIG_FILE -e
    " 2>&1
else
    msg_error "The program sshd is not installed. Install it by executing:\n pkg install openssh (in Termux app)"
    touch down
    sv stop .
    exit 1
fi
